# supervisord.conf
[supervisord]
autostart=true
autorestart=true
nodaemon=true

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true

[program:cron]
command=cron -f
autostart=true
autorestart=true
